# docker/supervisord.conf
# ─────────────────────────────────────────────────────────────────────────────
# HotelixOS — Supervisord du template établissement
# Emplacement cible : /etc/supervisor/conf.d/app.conf
#
# Gère deux processus, avec logs séparés :
#   - php-fpm (priorité 5  : démarre en premier)
#   - nginx   (priorité 10 : démarre après php-fpm)
# ─────────────────────────────────────────────────────────────────────────────

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
startsecs=2
startretries=5
stopsignal=QUIT
priority=5
stdout_logfile=/var/log/supervisor/php-fpm.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
priority=10
stdout_logfile=/var/log/supervisor/nginx.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/supervisor/nginx.err.log
stderr_logfile_maxbytes=0
